# deploy/booking-concierge-cleanup.timer
[Unit]
Description=Daily log purge for booking-concierge
Requires=booking-concierge-cleanup.service

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
